# gateway/gateway.service
[Unit]
Description=Agent Gateway - Hermes ↔ AI 通信桥梁
After=network.target

[Service]
Type=simple
User=admin
WorkingDirectory=/home/admin/hermes-platform/gateway
EnvironmentFile=/etc/systemd/system/agent-gateway.service.d/env.conf
ExecStart=/home/admin/praison-env/bin/python /home/admin/hermes-platform/gateway/gateway.py
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
